# attractors/python/docker/attractors.cron
@daily /opt/attractors/web/run.sh
# Intentionally empty line to make this a valid cron file
